# meta-nanopi-zero2/recipes-core/rauc-handlers/files/rauc-mark-good.service
[Unit]
Description=RAUC Mark Good - Validate successful boot
Documentation=https://rauc.io/
After=multi-user.target docker.service rauc.service rauc-boot-check.service
Wants=docker.service
RequiresMountsFor=/persist

[Service]
Type=oneshot
ExecStart=/usr/lib/rauc/rauc-mark-good.sh
RemainAfterExit=yes
TimeoutStartSec=120
StandardOutput=journal
StandardError=journal
# Give services time to start before validation
ExecStartPre=/bin/sleep 10

[Install]
WantedBy=multi-user.target
