# meta-jsdelivr/recipes-jsdelivr/jsdelivr-service/files/jsdelivr-docker-controller.service
[Unit]
Description= JSDELIVR docker controller service
After=network-online.target docker.socket firewalld.service
Wants=network-online.target
Requires=docker.socket

[Service]
# Type=simple, not notify: the python controller does not call sd_notify(READY=1),
# so a notify type leaves the start job hung forever and blocks multi-user.target.
Type=simple
ExecStart=python3 /usr/bin/jsdelivr-docker-controller.py
ExecReload=/bin/kill -s HUP $MAINPID

# kill only the docker process, not all processes in the cgroup
KillMode=process
TimeoutSec=infinity


[Install]
WantedBy=multi-user.target
